# cloudnest.service
[Unit]
Description=CloudNest Backend Service
After=network.target mongodb.service
Wants=mongodb.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/cloudnest
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# Load environment variables from .env file
EnvironmentFile=/var/www/cloudnest/.env

# Start the application
ExecStart=/usr/bin/python3 /var/www/cloudnest/app.py

# Restart policy
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudnest

[Install]
WantedBy=multi-user.target
